# services/listen_sensible_data/listen_sensible_data.socket
[Unit]
Description=Listen Sensible Data Socket
PartOf=listen_data.service

[Socket]
ListenStream=0.0.0.0:8886

[Install]
WantedBy=sockets.target